# etc/supervisor/conf.d/kibana.conf
[program:kibana]
command=/opt/kibana/bin/kibana
