# etc/supervisor/conf.d/kibana.conf
[program:kibana]
command=/usr/share/kibana/bin/kibana
